# packaging/stealthwg.service
[Unit]
Description=StealthWG masked WireGuard server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/stealthwg up
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
